# deployment/systemd/discord-red-roller-watcher.path
[Path]
PathModified=/etc/discord-red-roller

[Install]
WantedBy=multi-user.target
